# contrib/podman/dragonrouter-redis.container
[Unit]
Description=Dragon Router Redis (rate-limiter backend)
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/redis:7-alpine
ContainerName=dragonrouter-redis

Network=dragonrouter.network
NetworkAlias=redis

Volume=dragonrouter-redis.volume:/data

Exec=redis-server --save 60 1 --loglevel warning

HealthCmd=/bin/sh -c "redis-cli ping | grep -q PONG"
HealthInterval=10s
HealthRetries=3
HealthTimeout=5s

[Service]
Restart=always

[Install]
WantedBy=default.target
